# smartapplianceenabler/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Smart Appliance Enabler
# Configures NGINX for use with the Smart Appliance Enabler
# ==============================================================================


bashio::var.json \
    supervisor_token "${SUPERVISOR_TOKEN}" \
    | tempio \
        -template /etc/nginx/nginx.conf.template \
        -out /etc/nginx/nginx.conf
