# smartapplianceenabler/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Smart Appliance Enabler
# Configures NGINX for use with the Smart Appliance Enabler
# ==============================================================================

echo $SUPERVISOR_TOKEN}
envsubst "${SUPERVISOR_TOKEN}" < /etc/nginx/nginx.conf.template > /etc/nginx/nginx.conf
