# smartapplianceenabler/rootfs/etc/s6-overlay/s6-rc.d/init-smartapplianceenabler/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Smart Appliance Enabler
# Configure the Smart Appliance Enabler
# ==============================================================================

bashio::var.json \
    semp_port "$(bashio::addon.port 55555)" \
    server_port "$(bashio::addon.port 8080)" \
    | tempio \
        -template /etc/default/smartapplianceenabler.template \
        -out /etc/default/smartapplianceenabler

