# smartapplianceenabler/rootfs/etc/s6-overlay/s6-rc.d/init-smartapplianceenabler/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Smart Appliance Enabler
# Configure the Smart Appliance Enabler
# ==============================================================================

SEMP_PORT="$(bashio::addon.port 55555)"
echo $SEMP_PORT
envsubst "${SEMP_PORT}" < /etc/default/smartapplianceenabler.template > /etc/default/smartapplianceenabler
