# smartapplianceenabler/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Smart Appliance Enabler
# Take down the S6 supervision tree when Smart Appliance Enabler fails
# ==============================================================================

exec /run/s6/basedir/bin/halt
